# bridges/slack/launchd.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Example launchd plist for the Slack bridge on macOS.

  Replace every {{REPO_ROOT}} with the absolute path you cloned
  @elnora-ai/linear into (or wherever bridges/slack/bridge.py lives on
  your machine), then install via:

      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.example.linear-curator-bridge.plist

  Fires at 09:38 + 11:30 + 14:30 on weekdays — `post-pending` first, then
  `resolve` (the `tick` mode). Adjust StartCalendarInterval to your cadence.

  Important: launchd does not source .envrc, .env, or shell rc files.
  Either set SLACK_BOT_TOKEN / ANTHROPIC_API_KEY in EnvironmentVariables
  below (mind the file permissions — plists live under ~/Library/LaunchAgents
  which is mode 0700 by default), or wrap the python invocation in a small
  script that loads your secret store before exec'ing bridge.py.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.linear-curator-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-u</string>
        <string>{{REPO_ROOT}}/bridges/slack/bridge.py</string>
        <string>tick</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>38</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>38</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>38</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>38</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>38</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>{{REPO_ROOT}}/bridges/slack/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>{{REPO_ROOT}}/bridges/slack/bridge-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <!--
          Uncomment + populate from your secret store. Do NOT commit a
          plist that contains real tokens.

          <key>SLACK_BOT_TOKEN</key>
          <string>xoxb-...</string>
          <key>ANTHROPIC_API_KEY</key>
          <string>sk-ant-...</string>
        -->
    </dict>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}</string>
</dict>
</plist>
